# docker/files/etc/services.d/1.s3upload/run
#!/usr/bin/with-contenv sh
pure-uploadscript -B -r /opt/uploadscript.sh
